# config/services/balena-registry.service
[Unit]
Description=balena-registry
Requires=confd.service
After=confd.service
StartLimitIntervalSec=0

[Service]
StandardOutput=journal+console
StandardError=journal+console
WorkingDirectory=/usr/src/app
EnvironmentFile=/usr/src/app/config/env
ExecStart=/usr/src/app/entry.sh
Restart=always
SyslogIdentifier=registry

[Install]
WantedBy=basic.target
