# config/services/redis-server.service.d/restart.conf
[Unit]
# Always retry starting redis in case it dies for some reason
StartLimitIntervalSec=0
